# userjob/userjob.timer
[Unit]
Description = User job

[Timer]
Persistent = true
OnCalendar = *-*-* *:00:00

[Install]
WantedBy = timers.target
