# userjob/userjob.service
[Unit]
Description = User job

[Service]
ExecStart = /bin/sh -c '${HOME}/.config/userjob/exec'
Type = oneshot
